# docker/supervisor/gunicorn.ini
[program:gunicorn]
user=root
priority=900
command=gunicorn wsgi:app -b 0.0.0.0:8000
autostart=true
startsecs=3
autorestart=unexpected
directory=/usr/app
stdout_logfile=/var/log/gunicorn_stdout.log
redirect_stderr=false
stderr_logfile=/var/log/gunicorn_stderr.log
stopsignal=QUIT
