# docker/supervisor/nginx.ini
[program:nginx]
user=root
priority=900
command=nginx -g 'daemon off;'
autostart=true
startsecs=3
autorestart=unexpected
directory=/tmp
stdout_logfile=/var/log/nginx_stdout.log
redirect_stderr=false
stderr_logfile=/var/log/nginx_stderr.log
stopsignal=QUIT
